# hardware/95-monitor-hotplug.rules
ACTION=="change", SUBSYSTEM=="drm", RUN+="~/.config/i3/monitor-handler.sh"